# xdg-autostart.conf.in
description "A process defined by an XDG Autostart file"

start on xdg-autostart-start
stop on desktop-end

instance $XDG_AUTOSTART_NAME

pre-start exec @pkglibexecdir@/xdg-autostart-parse
script
	echo $XDG_AUTOSTART_NAME
#	@pkglibexecdir@/exec-line-exec
end script
