# xdg-autostart-find.conf.in
description "Look through the XDG data directories to find things to autostart"

start on desktop-start
task

emits xdg-autostart-start

exec @pkglibexecdir@/xdg-autostart-find
